# etc/periodic/hourly/000.zfs-auto-snapshot
#!/bin/sh

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

case "${hourly_zfs_snapshot_enable:-yes}" in
    [Yy][Ee][Ss])
        for pool in ${hourly_zfs_snapshot_pools:-//}
        do
            /sbin/zfs-auto-snapshot -q -g -l hourly -k "${hourly_zfs_snapshot_keep:-24}" "$pool"
        done
        ;;
    *)
        ;;
esac
