# etc/periodic/monthly/000.zfs-auto-snapshot
#!/bin/sh

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

case "${monthly_zfs_snapshot_enable:-yes}" in
    [Yy][Ee][Ss])
        for pool in ${monthly_zfs_snapshot_pools:-//}
        do
            /sbin/zfs-auto-snapshot -q -g -l monthly -k "${monthly_zfs_snapshot_keep:-12}" "$pool"
        done
        ;;
    *)
        ;;
esac
